# laravel-queue.service
[Unit]
Description=Laravel Queue Worker for Broadcasting
After=network.target redis.service
Wants=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/your-project
ExecStart=/usr/bin/php artisan queue:work sync --sleep=3 --tries=3 --max-time=3600 --timeout=60
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=laravel-queue

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/your-project/storage

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment variables
Environment=APP_ENV=production
Environment=APP_DEBUG=false

[Install]
WantedBy=multi-user.target
